# bash/.config/bash/profile.d/history.sh
export HISTCONTROL=ignoredups
export HISTSIZE=10000
export HISTTIMEFORMAT='%Y%m%d%H%M%S '
