# services/.config/systemd/user/xinit.target
[Unit]
Description=Start an X Session
Requires=environment.target
IgnoreOnIsolate=true

[Install]
Alias=default.target
